# conan.container
[Unit]
Description=CONAN — Gestion des concerts Négitachi
After=network-online.target
Wants=network-online.target

[Container]
Image=conan:latest
PublishPort=8080:80

# Volume persistant pour les données
Volume=/opt/conan/data:/var/www/html/data:Z

# Variables d'environnement
Environment=APACHE_RUN_USER=www-data
Environment=APACHE_RUN_GROUP=www-data

# Sécurité Podman
NoNewPrivileges=true
ReadOnly=false

[Service]
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
